# scriptkit/launchd/com.scriptkit.diskcleaner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.scriptkit.diskcleaner</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/YOUR_USERNAME/scriptkit/scripts/run_clean_disk.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
      <key>Weekday</key>
      <integer>6</integer> <!-- Friday -->
      <key>Hour</key>
      <integer>23</integer> <!-- 11 PM -->
      <key>Minute</key>
      <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/clean_disk.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/clean_disk.err</string>

    <key>RunAtLoad</key>
    <true/>
  </dict>
</plist>
